# upxtreme_v2/193/systems/multi-user.target.wants/snap-snap\x2dstore-599.mount
[Unit]
Description=Mount unit for snap-store, revision 599
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/snap-store_599.snap
Where=/snap/snap-store/599
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=default.target multi-user.target
